# deployment/drakon-node.service
[Unit]
Description=Drakon Blockchain Node
After=network.target

[Service]
Type=simple
User=drakon
Group=drakon
WorkingDirectory=/opt/drakon-node
Environment=NODE_ENV=production
Environment=HOME=/home/drakon
ExecStart=/usr/bin/npm run start:prod
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=drakon-node

[Install]
WantedBy=multi-user.target 